# bin/psistats-service
#!/bin/sh
### BEGIN INIT INFO
# Provides:          psistats
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      1
# Short-Description: psistats linux client
# Description:       psistats linux client
### END INIT INFO

# Author: v0idnull <v0idnull@psikon.com>

# Do NOT "set -e"

PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
DESC="Psistats Linux Client"
HOME=/opt/psistats
NAME=psistats
PIDFILE=/var/run/$NAME.pid
SCRIPT=$NAME

# Read configuration variable file if it is present
[ -r /etc/default/$NAME ] && . /etc/default/$NAME

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

case "$1" in
    start)
        $SCRIPT start
        ;;
    stop)
        $SCRIPT stop
        ;;
    reload|restart)
        $SCRIPT stop
        $SCRIPT start
        ;;
    status)
        $SCRIPT status
        ;;
    *)
        $SCRIPT
        exit 1
esac
exit 0

